# pkg/debian/DEBIAN/postinst
ldconfig
! test -e /etc/c8ydemo.conf && touch /etc/c8ydemo.conf
systemctl enable c8ydemo
