# pkg/debian/DEBIAN/postinst
ldconfig
! test -e /etc/$AGENT_NAME.conf && touch /etc/$AGENT_NAME.conf
systemctl enable $AGENT_NAME
