# pkg/debian/DEBIAN/prerm
systemctl disable cumulocity-agent
systemctl disable cumulocity-remoteaccess